# schedule_api_gspread/etc/crontab/run_splatoon_uploader
#!/bin/bash

cd /opt/splatoon_venv
source ./bin/activate

cd PATH
/opt/splatoon_venv/bin/python3 main.py

# 8 0 * * * ..../etc/crontab/run_splatoon_uploade
